# cul2mqtt/rootfs/etc/services.d/cul2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

exec cul-mqtt \
    --device "$(bashio::config 'device')" \
    --broker "$(bashio::services mqtt 'host')" \
    --topic "$(bashio::config 'topic')" \
    --username "$(bashio::services mqtt 'username')" \
    --password "$(bashio::services mqtt 'password')"
